# contrib/kb-sentinel.service
[Unit]
Description=KB Sentinel - 4by4 Keyboard MQTT Bridge
Documentation=https://github.com/saikhurana98/kb_sentinel
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/kb_sentinel
ExecStart=%h/kb_sentinel/.venv/bin/python3 %h/kb_sentinel/kb_sentinel.py
Restart=always
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=10
StandardOutput=journal
StandardError=journal
# Restart policy
StartLimitInterval=300
StartLimitBurst=3

[Install]
WantedBy=default.target